# ansible/templates/vepc.service.j2
[Unit]
Description=vEPC - Evolved Packet Core Emulator
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/cubasnik/vEPC

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{ vepc_repo_path }}
Environment=CONFIG_PATH={{ config_path }}
Environment=VEPC_TRAFFIC_PORTS={{ traffic_linux_ports | join(',') }}
Environment=HOST_SYS_CLASS_NET=/sys/class/net
UMask=0000

ExecStartPre=/usr/bin/mkdir -p {{ log_path }} {{ config_path }}
ExecStartPre=/usr/bin/rm -f /tmp/vepc.sock
ExecStart={{ build_path }}/vepc --config={{ config_path }}
ExecStopPost=/usr/bin/rm -f /tmp/vepc.sock

Restart=on-failure
RestartSec=5
LimitNOFILE=65535
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=false
ProtectHome=no

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vepc

[Install]
WantedBy=multi-user.target
